# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server /etc/redis/redis.conf
autostart=true
startsecs=0

[program:new_api]
command=newapi
autostart=true
startsecs=0